# dot-config/systemd/user/vim-plugins-update.timer
[Unit]
Description=Update vim plugins daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
