# files/vim-plugins-update.service
[Unit]
Description=Update vim plugins
Wants=vim-ycm-update.service
[Service]
Type=oneshot
ExecStart=vim +PluginUpdate +qall
