# files/zoom-update.timer
[Unit]
Description=Run Zoom update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
